# systemd/nubeio-enable-uart-pins.service
[Unit]
Description=Enable UART pins

[Service]
Type=simple
ExecStart=/home/debian/bonescript-api/systemd/nubeio-enable-uart-pins.sh
